# sbrDocumentApi.service
[Service]
ExecStart=/usr/share/sbr/sbr-document-api/start.sh
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=sbr-document-api
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
